# 2025-final/web/customsconundrum/container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid

[program:proxy]
command=python3 main.py
directory=/proxy
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/proxy_err.log
stdout_logfile=/var/log/supervisor/proxy_out.log
environment=PYTHONUNBUFFERED=1

[program:app]
command=uvicorn main:app --host 0.0.0.0 --port 8081
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/app_err.log
stdout_logfile=/var/log/supervisor/app_out.log
environment=PYTHONUNBUFFERED=1

[program:admin]
command=uvicorn main:app --host 0.0.0.0 --port 8082
directory=/admin
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/admin_err.log
stdout_logfile=/var/log/supervisor/admin_out.log
environment=PYTHONUNBUFFERED=1

[program:bot]
command=uvicorn main:app --host 0.0.0.0 --port 8083
directory=/bot
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/bot_err.log
stdout_logfile=/var/log/supervisor/bot_out.log
environment=PYTHONUNBUFFERED=1